# salt/roots/ice-registry-tpl.conf
description "IceGrid Registry Service"

{% if pillar['userid'] == 'vagrant' %}
start on vagrant-mounted MOUNTPOINT=/vagrant
{% else %}
start on runlevel [2]
{% endif %}
stop on runlevel [!2]

console log

setuid {{ pillar['userid'] }}

exec icegridregistry --Ice.Config={{ pillar['app_root'] }}/grid/registry.cfg
